# services/pegabot/run
#!/bin/bash
export USER=node

su - node -c 'cd /src; npm start'
